# tinydns/trunk/tinydns.openrc
#!/sbin/openrc-run
command=/usr/bin/tinydns
command_args="-d"
pidfile=
name="tinydns"
description="tiny DNS cache service"
